# qpacman_service/qpacman.service.in
[Unit]
Description=QPacman DBus service

[Service]
Type=dbus
BusName=com.alexl.qt.QPacmanService
ExecStart=$$INSTALL_PREFIX/bin/qpacman_service
Environment="PATH=$$INSTALL_PREFIX/bin:$$INSTALL_PREFIX/sbin"

[Install]
WantedBy=multi-user.target
